# cron.d/db_limit
*/2 * * * * /usr/bin/php /etc/dbmax.php
