# cron.d/db_limit
*/10 * * * * /usr/bin/php /etc/scripts/dbmax.php
