# cron.d/incron_refresh
*/15 * * * * root /etc/scripts/incron_config.sh
