# cron.d/remote_backup
00 13 * * * root /etc/scripts/remote_backup.sh
